# tools/init/runit/zorix/run
#!/bin/sh
exec chpst -u zorix /usr/local/bin/zorix -c /etc/zorix/config.toml
